# transmission.conf
[program:transmission]
autostart=true
autorestart=true
command=/usr/bin/transmission-daemon  -f -g /var/lib/transmission-daemon/info/
process=(program_name)s
